# docker/root/etc/services.d/budgeteer/run
#!/usr/bin/with-contenv bash

 exec budgeteer --docker --log-level=$LOGLEVEL
